# config/systemd-user/experiments.slice
[Unit]
Description=ML Experiments slice
Before=slices.target

[Slice]
# Values templated by deploy.sh from config/resources.conf
CPUQuota=200%
MemoryMax=24G
MemoryHigh=20G
MemorySwapMax=0
IOWeight=50
ManagedOOMMemoryPressure=kill
