# config/systemd-user/reset-failed.timer
[Unit]
Description=Hourly cleanup of failed systemd-run units

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
